# hub_cron/rootfs/etc/cron.d/common_cron
*/5 * * * * root python3 /data/commands/auto_run.py >> /var/log/syslog